# infra/demo-cluster/pln.service
[Unit]
Description=PLN Mesh Node
After=network-online.target

[Service]
Type=simple
User=pln
Group=pln
# PLN_DEBUG_PPROF exposes /debug/pprof on the metrics HTTP server. Safe for
# the SG-restricted demo; should not be set in production.
Environment=PLN_DEBUG_PPROF=1
ExecStart=/usr/bin/pln --dir /var/lib/pln up --control-addr :50051 --http :9090
Restart=on-failure
RestartSec=5
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
