# infra/prod/laptop/sh.pln.prod.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!--
  Launch agent for the laptop-side prod pln daemon. Runs alongside the
  brew-managed home-cluster daemon on a different port and state dir.

  Install via `just prod-laptop-install` (which templates __HOME__ to
  the user's home path at install time). Logs land at
  ~/.pln-prod/{stdout,stderr}.log. The daemon is the root of trust for
  the prod mesh — keep the state dir backed up.
-->
<dict>
    <key>Label</key>
    <string>sh.pln.prod</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/pln</string>
        <string>--dir</string>
        <string>__HOME__/.pln-prod</string>
        <string>up</string>
        <string>--port</string>
        <string>60711</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/.pln-prod/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/.pln-prod/stderr.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
