# infra/demo-cluster/pln.service.d/demo.conf
# Demo-only overrides layered on top of the unit `pln service install`
# writes. Control API and Prometheus bind addresses are set via
# `pln set control-addr :50051 && pln set http :9090` during provisioning.
#
# MemoryDenyWriteExecute is relaxed so wazero can use its JIT compiler
# instead of falling back to the interpreter — demo load tests care
# about throughput.

[Service]
Environment=PLN_DEBUG_PPROF=1
MemoryDenyWriteExecute=no
